# engineering-calculation-system-skill-pack-v2.3.0/core/engineering-calculation-system/project_template/engineering_calc_project/deploy/systemd/engineering-calc.service
[Unit]
Description=Engineering Calculation Web Application
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/engineering-calc
EnvironmentFile=/opt/engineering-calc/deploy/env.example
Environment=PYTHONPATH=/opt/engineering-calc/src
ExecStart=/opt/engineering-calc/.venv/bin/gunicorn "webapp.app:create_app()" --bind 127.0.0.1:5000 --workers 2
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
